# chb3.service
[Unit]
Description=ChronophylosBot Version 3

[Service]
ExecStart=chb3

[Install]
WantedBy=multi-user.target
